# geoipupdate3.cron
#!/bin/sh

# This updates the IPv4 databases based on product numbers in /etc/GeoIP.conf
geoipupdate3 > /dev/null
